# src/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/bin/bash -c "exec /usr/sbin/httpd start"

[program:httpd]
command=/bin/bash -c "source /etc/httpd/envvars && exec /usr/sbin/httpd -D FOREGROUND"
